# vscode-remote/rootfs/etc/profile.d/homeassistant.sh
export HASS_SERVER="http://supervisor/core"
export HASS_TOKEN="${SUPERVISOR_TOKEN:-}"
